# etc/systemd/system/multi-level-test.target
[Unit]
Description=multi-level-test service
Wants=multi-level-test-a.target
Wants=multi-level-test-b.target

[Install]
WantedBy=multi-user.target
